# launchd/bareclaude.price-pulse.daily-alert.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>bareclaude.price-pulse.daily-alert</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/daelegbe/BareClaude/finn/repos/price-pulse/scripts/price-alert-run.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DB_PATH</key>
    <string>/Users/daelegbe/BareClaude/finn/.state/price-pulse.db</string>
    <key>TELEGRAM_SEND_SCRIPT</key>
    <string>/Users/daelegbe/BareClaude/finn/scripts/telegram-send.sh</string>
  </dict>

  <!-- Daily at 09:00 Mountain Time (launchd uses system local clock) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/daelegbe/BareClaude/finn/.state/price-alert-run.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/daelegbe/BareClaude/finn/.state/price-alert-run.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
